# unit-files/check_health.timer
# bitcore.timer
[Unit]
Description=Wait 5 min to start bitcore on boot

[Timer]
OnBootSec=5min
[Install]
WantedBy=multi-user.target
