# unit-files/check_ips.timer
# check_ips.timer
[Unit]

Description=Check IPs for haproxy Timer

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
