# unit-files/haproxy.service
# haproxy.service
[Unit]
Description=HAproxy service for btc
After=docker.service
Requires=docker.service
After=create_docker_network.service
Requires=create_docker_network.service

[Service]
TimeoutStartSec=5
Restart=on-failure

ExecStartPre=/bin/sleep 1
ExecStartPre=-/usr/bin/docker stop bitcoin_haproxy
ExecStartPre=-/usr/bin/docker rm bitcoin_haproxy
ExecStartPre=/usr/bin/docker pull haproxy:alpine
ExecStart=/usr/bin/docker run \
  --net=bitcoind \
  -m 128m \
  -p 8332:8332 \
  -p 8333:8333 \
  -p 1024:1024 \
  --expose=8332 \
  --expose=8333 \
  --expose=1024 \
  -v /bitcoind/configs/haproxy/haproxy.cfg:/usr/local/etc/haproxy/haproxy.cfg \
  --name bitcoin_haproxy \
  haproxy:alpine \
/usr/local/sbin/haproxy -f /usr/local/etc/haproxy/haproxy.cfg -L %H

ExecStop=/usr/bin/docker stop bitcoin_haproxy
ExecStopPost=/usr/bin/docker rm -f bitcoin_haproxy
ExecReload=/usr/bin/docker restart bitcoin_haproxy


[Install]
WantedBy=multi-user.target
